# conf/api-uwsgi.conf
start on [2345]
stop on [06]

script
    cd /home/bhs/api
    exec uwsgi --ini /etc/flask-uwsgi/flask-uwsgi.ini
end script
